# systemd/chameleon-webui.service
[Unit]
Description=Chameleon WebUI Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/chameleon-webui
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
